# compilers/extra/c3/update.cron
0 * * * * /update/c3 &> /cron-log/c3.log
